# mm/pytools/grinsdb/genlicense.crontab
30 3 * * 0	/ufs/jack/src/mm/pytools/grinsdb/genlicense.sh
